# centos-minimal/etc/profile.d/java-env.sh
#!/bin/bash
export JAVA_HOME="/usr/java/jdk1.8.0_74/bin/java"
export JRE_HOME="/usr/java/jdk1.8.0_74/jre/bin/java"
PATH=$PATH:$HOME/bin:JAVA_HOME:JRE_HOME

